# deploy/themain-backup.timer
[Unit]
Description=THE MAIN — запуск ежедневного бэкапа базы

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
